# proc/supervisord.conf
[supervisord]
nodaemon=true

[program:sacosapi]
command=/app/main
autostart=true
autorestart=true
redirect_stderr=true
